# maxscale-cdc-adapter/test/docker/mcs/service/systemd-journald/run
#!/bin/sh
/usr/lib/systemd/systemd-journald 2>&1
