# kafka-avro-adapter/docker/mcs/service/columnstore/run
#!/bin/sh
exec 2>&1
(mkdir -p /var/lock/subsys && touch /var/lock/subsys/columnstore) 
exec ProcMon > /dev/null
